# xremap/xremap.service
[Unit]
Description=xremap

[Service]
Type=simple
Restart=always
ExecStart=xremap --watch=config %h/.config/xremap/wordstar.yml

[Install]
WantedBy=default.target
